# ErmisServer/Installers/Unix/DEB/ermis-server-installer_all/etc/systemd/system/ermis-server.service
[Unit]
Description=Ermis-Server Service
After=network.target

[Service]
User=ermis
Type=simple

LoadCredential=keystore_password:/etc/ermis-server/keystore_password.cred
LoadCredential=db_user_password:/etc/ermis-server/db_user_password.cred
LoadCredential=emailer_password:/etc/ermis-server/emailer_password.cred

WorkingDirectory=/opt/ermis-server/
ExecStart=/bin/bash -c 'cd /opt/ermis-server/ && ./ermis-server.sh'
Restart=no

[Install]
WantedBy=multi-user.target
